# dot_config/systemd/user/executable_tailscale-systray.service
[Unit]
Description=Tailscale System Tray
# graphical.target is a system target; user tray needs the user graphical
# session (niri) to be up before the tray host exists.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/tailscale systray

[Install]
# Activated via niri.service.wants/ (chezmoi-managed symlink).
WantedBy=niri.service
