# dot_config/systemd/user/handy.service
[Unit]
Description=Handy speech-to-text
Documentation=https://github.com/cjpais/handy
After=graphical-session.target pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do ls /run/user/$(id -u)/wayland-* >/dev/null 2>&1 && exit 0; sleep 0.5; done; echo "Warning: Wayland socket not found after 15s"'
ExecStart=/usr/bin/handy --start-hidden --no-tray
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
