# dot_config/systemd/user/update-flatpak-package-files.service.tmpl
[Unit]
Description=Update flatpak declarative package list

[Service]
Type=oneshot
ExecStart=/bin/bash {{ .chezmoi.homeDir }}/flatpak/update-package-files.sh
