# dot_config/systemd/user/update-flatpak-package-files.path
[Unit]
Description=Watch for flatpak app changes

[Path]
PathChanged=/var/lib/flatpak/app
TriggerLimitIntervalSec=0

[Install]
WantedBy=default.target
