# dot_config/systemd/user/xremap-device-watchdog.timer
[Unit]
Description=Poll Go60/kanata input devices for xremap restarts

[Timer]
OnBootSec=20s
OnUnitActiveSec=10s
AccuracySec=2s
Unit=xremap-device-watchdog.service

[Install]
WantedBy=timers.target
