# dot_config/systemd/user/xremap-device-watchdog.service
[Unit]
Description=Restart xremap when Go60/kanata input devices change

[Service]
Type=oneshot
ExecStart=%h/.config/xremap/restart-if-devices-changed.sh
