# dot_config/scarlett-kvm/scarlett-kvm-reconnect.service
[Unit]
Description=Reconnect Scarlett 8i6 audio after KVM switch

[Service]
Type=oneshot

# Wait for the KVM hub to settle — it bounces multiple times during a switch.
# Starting audio playback during a bounce causes snd_usb_audio URB timeouts
# that deadlock the xHCI controller, killing ALL USB (hard reboot required).
ExecStartPre=/bin/sleep 10

# Verify the Scarlett is still present after settling.
# If the KVM switched away during the wait, exit cleanly.
ExecStartPre=/bin/bash -c 'grep -rqs "Scarlett" /sys/bus/usb/devices/*/product'

# Recreate virtual sinks and loopbacks (also restores ALSA mixer state).
ExecStart=/usr/bin/systemctl --user --machine=mg@ restart system-audio-sink.service
